# srv/projects/kb/data/supervisor/supervisord.conf
[supervisord]
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid
nodaemon=true

[program:qdrant]
command=/usr/local/bin/qdrant
directory=/data/qdrant
autostart=true
autorestart=true
startsecs=5
priority=5
stdout_logfile=/var/log/qdrant/stdout.log
stderr_logfile=/var/log/qdrant/stderr.log

[program:ollama]
command=/usr/local/bin/ollama serve
directory=/data/ollama
autostart=true
autorestart=true
startsecs=5
priority=8
stdout_logfile=/var/log/ollama/stdout.log
stderr_logfile=/var/log/ollama/stderr.log
environment=OLLAMA_MODELS="/data/ollama"

[program:uvicorn]
command=uvicorn app.main:app --host 0.0.0.0 --port 8000
redirect_stderr=true
autostart=true
autorestart=true
priority=10

[program:nginx]
command=/usr/sbin/nginx -c /data/nginx.conf -g 'daemon off;'
redirect_stderr=true
autostart=true
autorestart=true
priority=20

[program:watchdog]
command=/bin/bash /data/scripts/watchdog.sh
startsecs=5
autostart=true
autorestart=true
exitcodes=0,2
priority=30
