# srv/projects/kb/data/supervisor/supervisord.conf
[supervisord]
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid
nodaemon=true

[program:uvicorn]
command=uvicorn app.main:app --host 0.0.0.0 --port 8000
redirect_stderr=true
autostart=true
autorestart=true
priority=10

[program:nginx]
command=/usr/sbin/nginx -c /data/nginx.conf -g 'daemon off;'
redirect_stderr=true
autostart=true
autorestart=true
priority=20

[program:watchdog]
command=/bin/bash /data/scripts/watchdog.sh
startsecs=5
autostart=true
autorestart=true
exitcodes=0,2
priority=30
