# app/supervisord.conf
[supervisord]
nodaemon=true
logfile=/srv/projects/kb/data/logs/supervisor/supervisord.log
logfile_maxbytes=50MB
pidfile=/srv/projects/kb/data/supervisor/supervisord.pid

[program:ollama]
command=/usr/local/bin/ollama serve
directory=/srv/projects/kb/data/storage/ollama
environment=OLLAMA_HOST=0.0.0.0:11434,OLLAMA_MODELS=/srv/projects/kb/data/storage/ollama,HOME=/srv/projects/kb
user=knowlab
autorestart=true
startretries=10
stdout_logfile=/srv/projects/kb/data/logs/ollama/stdout.log
stderr_logfile=/srv/projects/kb/data/logs/ollama/stderr.log
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:qdrant]
command=/usr/local/bin/qdrant
directory=/srv/projects/kb/data
environment=QDRANT__STORAGE__STORAGE_PATH=/srv/projects/kb/data/storage/qdrant,HOME=/srv/projects/kb
user=knowlab
autorestart=true
startretries=10
stdout_logfile=/srv/projects/kb/data/logs/qdrant/stdout.log
stderr_logfile=/srv/projects/kb/data/logs/qdrant/stderr.log
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:api]
command=/usr/local/bin/uvicorn app.main:app --host 127.0.0.1 --port 8080 --proxy-headers --log-level=info
directory=/srv/projects/kb/app
user=knowlab
autorestart=true
startretries=10
stdout_logfile=/srv/projects/kb/data/logs/api/stdout.log
stderr_logfile=/srv/projects/kb/data/logs/api/stderr.log
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
autorestart=true
startretries=10
stdout_logfile=/srv/projects/kb/data/logs/nginx/stdout.log
stderr_logfile=/srv/projects/kb/data/logs/nginx/error.log
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
